# build/crontab
*/15 * * * * * * app_cli app:nowplaying
0 */15 * * * * * app_cli app:update
