# build/crontab
SHELL=/bin/bash
BASH_ENV=/container.env
* * * * *   sleep 10; cd /var/app/www && php bin/console app:nowplaying
* * * * *   sleep 25; cd /var/app/www && php bin/console app:nowplaying
* * * * *   sleep 40; cd /var/app/www && php bin/console app:nowplaying
* * * * *   sleep 55; cd /var/app/www && php bin/console app:nowplaying
*/15 * * * *   cd /var/app/www && php bin/console app:update
