# build/crontab
*/15 * * * * * * php /var/app/www/bin/console app:nowplaying
0 */15 * * * * * php /var/app/www/bin/console app:update
