# runit/php-fpm/run
#!/bin/bash

/usr/sbin/php-fpm7.2 -F --fpm-config /etc/php/7.2/fpm/php-fpm.conf -c /etc/php/7.2/fpm/